# systemd/artwall.service
[Unit]
Description=Set a random Met Museum painting as the Sway wallpaper
Documentation=https://metmuseum.github.io/
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=oneshot
# __REPO__ and __PYTHON__ are filled in by install.sh to point at this checkout.
Environment=PYTHONPATH=__REPO__
ExecStart=__PYTHON__ -m artwall
